# setup/proxy/sampleconf/services.d/shibd/run
#!/usr/bin/execlineb -P
/usr/sbin/shibd -f -F
