# setup/proxy/sampleconf/services.d/httpd/run
#!/bin/sh
#/usr/bin/execlineb -P
rm /var/run/httpd.pid
/usr/sbin/httpd -DFOREGROUND
