# setup/proxy/sampleconf/supervisord.d/shibd.ini
[program:shibd]
command=/bin/bash -c "export LD_LIBRARY_PATH=/opt/shibboleth/lib64:${LD_LIBRARY_PATH} && /usr/sbin/shibd -f -F"

